# supervisord.conf
[supervisord]
nodaemon=true
 
[program:sshd]
command=/usr/sbin/sshd -D

[program:dashing]
#command=/bin/bash -c "exec dashing start"
